# docker/websocket_config/websockets.conf
[program:websockets]
command=/usr/bin/php /application/MyZoom/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
